# application/crontab
*/2 * * * * /usr/local/bin/python /response.py >> /var/log/cron.log 2>&1
